# videos-relay/supervisord.conf
[supervisord]
nodaemon=true

[program:relay]
command=/entrypoint.sh
autostart=true
autorestart=true